# supervisor.conf
# Supervisor Configuration for Payroll Management System
# Queue Management and Process Control

[program:payroll-system-queue]
command=php /var/www/payroll-system/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-queue.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-queue-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-horizon]
command=php /var/www/payroll-system/artisan horizon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-horizon.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-horizon-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-scheduler]
command=php /var/www/payroll-system/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-scheduler.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-scheduler-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-backup]
command=php /var/www/payroll-system/artisan backup:run
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-backup.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-backup-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-cleanup]
command=php /var/www/payroll-system/artisan data:cleanup
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-cleanup.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-cleanup-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-monitor]
command=php /var/www/payroll-system/artisan system:monitor
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-monitor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-monitor-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-export]
command=php /var/www/payroll-system/artisan export:process
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-export.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-export-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-import]
command=php /var/www/payroll-system/artisan import:process
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-import.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-import-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-notifications]
command=php /var/www/payroll-system/artisan notifications:send
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-notifications.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-notifications-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-reports]
command=php /var/www/payroll-system/artisan reports:generate
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-reports.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-reports-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-sync]
command=php /var/www/payroll-system/artisan integration:sync
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-sync.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-sync-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-webhooks]
command=php /var/www/payroll-system/artisan webhooks:process
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-webhooks.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-webhooks-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-audit]
command=php /var/www/payroll-system/artisan audit:process
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-audit.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-audit-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-archive]
command=php /var/www/payroll-system/artisan data:archive
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-archive.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-archive-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-restore]
command=php /var/www/payroll-system/artisan data:restore
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-restore.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-restore-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-health]
command=php /var/www/payroll-system/artisan health:check
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-health.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-health-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-performance]
command=php /var/www/payroll-system/artisan performance:monitor
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-performance.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-performance-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-security]
command=php /var/www/payroll-system/artisan security:scan
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-security.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-security-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-compliance]
command=php /var/www/payroll-system/artisan compliance:check
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-compliance.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-compliance-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-maintenance]
command=php /var/www/payroll-system/artisan maintenance:run
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-maintenance.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-maintenance-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-update]
command=php /var/www/payroll-system/artisan system:update
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-update.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-update-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-optimize]
command=php /var/www/payroll-system/artisan optimize
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-optimize.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-optimize-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-cache]
command=php /var/www/payroll-system/artisan cache:clear
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-cache.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-cache-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-config]
command=php /var/www/payroll-system/artisan config:cache
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-config.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-config-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-route]
command=php /var/www/payroll-system/artisan route:cache
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-route.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-route-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-view]
command=php /var/www/payroll-system/artisan view:cache
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-view.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-view-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-migrate]
command=php /var/www/payroll-system/artisan migrate --force
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-migrate.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-migrate-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-seed]
command=php /var/www/payroll-system/artisan db:seed --force
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-seed.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-seed-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-test]
command=php /var/www/payroll-system/artisan test
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-test.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-test-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

[program:payroll-system-deploy]
command=php /var/www/payroll-system/artisan deploy:run
autostart=false
autorestart=false
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/payroll-system-deploy.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/payroll-system-deploy-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=APP_ENV="production",APP_DEBUG="false"

# Group configuration for easier management
[group:payroll-system]
programs=payroll-system-queue,payroll-system-horizon,payroll-system-scheduler,payroll-system-monitor,payroll-system-notifications,payroll-system-sync,payroll-system-webhooks,payroll-system-audit,payroll-system-health,payroll-system-performance,payroll-system-security,payroll-system-compliance
priority=999

[group:payroll-system-maintenance]
programs=payroll-system-backup,payroll-system-cleanup,payroll-system-export,payroll-system-import,payroll-system-reports,payroll-system-archive,payroll-system-restore,payroll-system-maintenance,payroll-system-update,payroll-system-optimize,payroll-system-cache,payroll-system-config,payroll-system-route,payroll-system-view,payroll-system-migrate,payroll-system-seed,payroll-system-test,payroll-system-deploy
priority=1000

# Supervisor configuration
[supervisord]
nodaemon=false
minfds=1024
minprocs=200
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
identifier=supervisor
directory=/tmp
childlogdir=/var/log/supervisor
umask=022
serverurl=AUTO

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
username=admin
password=admin_password_here

[inet_http_server]
port=127.0.0.1:9001
username=admin
password=admin_password_here

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[eventlistener:memmon]
command=memmon -p payroll-system-queue=200MB -p payroll-system-horizon=300MB
events=TICK_60

[eventlistener:crashmail]
command=crashmail -m admin@your-domain.com
events=PROCESS_STATE_EXITED 